# udev/99-usb-sync.rules
ACTION=="add", SUBSYSTEM=="block", KERNEL=="sd[a-z][1-9]", RUN+="/usr/local/bin/usb_notify.sh"
